# mailredirector.service.template
[Unit]
Description=MailRedirector

[Service]
User={str}
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes
WorkingDirectory={str}

[Install]
WantedBy=multi-user.target